# luci-app-coova-chilli/root/etc/uci-defaults/99-coova-chilli
#!/bin/sh

uci rename chilli.@chilli[0]=settings
uci commit chilli
