# luci-app-coovachilli/root/etc/uci-defaults/99-coova-chilli
#!/bin/sh

uci rename chilli.@chilli[0]=settings
uci commit chilli

sleep 1

rm -f /etc/chilli/up.sh && mv /root/up.sh /etc/chilli/ && chmod +x /etc/chilli/up.sh
